# etc/supervisor-conf.d/supervisord-nginx.conf
[program:nginx]
command=start-nginx.sh
autorestart=true
autostart=true
